# base/root/etc/supervisor/conf.d/syslog-ng.conf
[program:syslog-ng]
command=/etc/supervisor/conf.d/syslog-ng
autorestart=true
